# supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=postgres
